# ml/devai-scheduler.service
[Unit]
Description=DevAI Fine-tuning Scheduler
After=network.target
Wants=network.target

[Service]
Type=simple
User=devai
Group=devai
WorkingDirectory=/opt/devai/ml
EnvironmentFile=/opt/devai/server/src/config/.env
Environment=PATH=/opt/devai/ml/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/devai/ml/venv/bin/python3 /opt/devai/ml/scripts/training_scheduler.py --min-pairs 200 --time 02:00
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
